# proton-wg@.service
[Unit]
Description=Proton WireGuard VPN for %i
After=network-online.target proton-killswitch.service
Wants=network-online.target
Requires=proton-killswitch.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=oneshot
EnvironmentFile=-/etc/proton/proton-common.env
EnvironmentFile=-/etc/proton/instances/%i/proton.env
EnvironmentFile=-/etc/proton/instances/%i/qbittorrent.env
ExecStart=/usr/local/bin/proton/proton-wg-up-safe.sh %i
ExecStop=/usr/local/bin/proton/proton-wg-down-safe.sh %i
RemainAfterExit=yes
TimeoutStartSec=120
TimeoutStopSec=60
UMask=0077
RuntimeDirectory=proton
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes
NoNewPrivileges=no
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
ProtectClock=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
ReadWritePaths=/run/proton /etc/wireguard/proton-runtime
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SETPCAP
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_SETPCAP

[Install]
WantedBy=multi-user.target
